# deploy/services/credit-spread.timer
[Unit]
Description=Start Bot 4 Credit Spread Seller on weekday mornings

[Timer]
# 9:55 AM ET — 5 min before entry window (10:00 AM)
OnCalendar=Mon..Fri 09:55:00
Persistent=false

[Install]
WantedBy=timers.target
