# deploy/services/credit-spread.service
[Unit]
Description=Bot 4 — 0DTE Credit Spread Seller
After=ibgateway.service
Requires=ibgateway.service

[Service]
Type=simple
User=optionsbot
Group=optionsbot
WorkingDirectory=/home/optionsbot/options_bot
Environment=DISPLAY=:1
EnvironmentFile=/home/optionsbot/options_bot/.env
ExecStart=/home/optionsbot/options_bot/.venv/bin/python credit_spread_live.py
Restart=on-failure
RestartSec=60
TimeoutStartSec=30

StandardOutput=append:/var/log/optbot/credit-spread.log
StandardError=append:/var/log/optbot/credit-spread.log

[Install]
WantedBy=multi-user.target
